# pkg/update-notifier/usr/lib/systemd/user/update-notifier.timer
[Unit]
Description=Daily check for new updates

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
